# docker/crontab
# Run DayZ Monitor every 2 hours
0 */2 * * * cd /app && python monitor.py >> /app/logs/cron.log 2>&1